# malakocut.service
[Unit]
Description=Malakocut NDR Sniffer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Ensure data directories exist
ExecStartPre=/usr/bin/mkdir -p /var/lib/malakocut/buffer /var/lib/malakocut/pcap
EnvironmentFile=/etc/default/malakocut
ExecStart=/usr/local/bin/malakocut \
    -interface "${MALAKOCUT_INTERFACE}" \
    -exporter "${MALAKOCUT_EXPORTER}" \
    -filter "${MALAKOCUT_FILTER}" \
    -blocklist "${MALAKOCUT_BLOCKLIST}"
Restart=always
RestartSec=5
StandardOutput=append:/var/log/malakocut.log
StandardError=append:/var/log/malakocut.log

[Install]
WantedBy=multi-user.target
